# packages/xmedia/camd/oscam/system.d/oscam.service
[Unit]
Description=Oscam Daemon
Requires=oscam-switch.service
After=oscam-switch.service network-online.target
Before=vdr.service
ConditionPathExists=/var/lock/start.oscam

[Service]
ExecStart=/bin/sh -c 'exec /usr/bin/oscam.start'
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
