# packages/sysutils/lcdproc/system.d/lcdd.service
[Unit]
Description=LCDProc
ConditionPathExists=/storage/.cache/services/lcdd.conf
After=multi-user.target

[Service]
ExecStart=/usr/bin/lcdproc
TimeoutStopSec=1s
Restart=always
RestartSec=3
StartLimitInterval=30
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
